# packaging/root_of_vm/home/pm/.config/autostart/pm_startup.desktop
[Desktop Entry]
Type=Application
Name=PMAutostart
Exec=/pm/scripts/autostart_after_x_login.sh
