# packaging/BaseDisk/pm_files/etc/systemd/system/openvpn_log.path
# Triggers output of the OpenVPN log file
[Unit]
Description=VPN log file output trigger
Requires=openvpn_pm.service

[Path]
PathChanged=/var/log/openvpn.log
